# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for superlists-staging.5e4n.com"

start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/sites/superlists-staging.5e4n.com/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/superlists-staging.5e4n.com.so$

